# openconnect-service/openconnect@.socket
[Unit]
Description=OpenConnect VPN socket - %i
Wants=network.target

[Socket]
ListenFIFO=/run/openconnect/%i
PipeSize=128

